# stage2/01-sys-tweaks/files/70-persistent-drm.rules
SUBSYSTEM=="platform", KERNEL=="*.v3d", ACTION=="add", RUN+="/usr/local/sbin/disable-drm-pm '%p'"
